# target/linux/rockchip/armv8/base-files/etc/hotplug.d/net/50-wifi-up
#!/bin/sh
wifi up
sleep 5
echo f > /sys/class/net/ra0/queues/rx-0/rps_cpus
echo f > /sys/class/net/rax0/queues/rx-0/rps_cpus
echo 10000 > /proc/sys/net/core/netdev_budget
